# misc/gutaur.timer
[Unit]
Description=Gutaur autopublish timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=default.target
