# misc/gutaur-autopublish.service
[Unit]
Description=Gutaur autopublish service
After=network.target

[Service]
ExecStart=/usr/bin/gutaur autopublish /app/aurs/*
