# misc/systemd/gutaur-autopublish@.service
[Unit]
Description=Gutaur autopublish service
After=network.target

[Service]
User=%i
ExecStart=/usr/bin/gutaur autopublish /app/aurs/*
